# scripts/com.forge-harness.frontier-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Replace every /path/to/... below with your own absolute paths before installing -->
    <key>Label</key>
    <string>com.forge-harness.frontier-digest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/path/to/forge-harness/scripts/frontier_digest_daily.sh</string>
    </array>

    <!-- Daily at 09:00 — change Hour to adjust -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/path/to/forge-harness/tracks/_meta/logs/launchd_out.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/forge-harness/tracks/_meta/logs/launchd_err.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/path/to/home/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/path/to/home</string>
    </dict>
</dict>
</plist>
